# docker/minecraftserver/docker.mcserver@.service
[Unit]
Description = Starts Minecraft Server Container @ Boot & Stops
After = docker.service
Requires = docker.service

[Service]
User = root
Group = root
Type = oneshot
# this pulls a file from the docker volume for the minecraft contiainer to
# get the ports below to open on container start. this file must exist. if you
# only intend to run a single world at a given time then adjust accordingly.
EnvironmentFile = /var/lib/docker/volumes/minecraftdata/_data/ENVFOLDER/%i
ExecStartPre = /usr/bin/docker run \
	-id \
	-h %i \
	-p ${MAINPORT}:25565 \
	--rm \
	--mount source=minecraft,destination=/var/lib/minecraft \
	--name %i \
	mcserver
ExecStart = /usr/bin/docker exec -i %i service minecraft start
ExecStart = /usr/bin/docker exec -i %i service cron start
ExecStop = /usr/bin/docker exec -i %i service minecraft stop
ExecStop = /usr/bin/docker stop %i
RemainAfterExit = yes

[Install]
WantedBy = multi-user.target
